# speech/transcribe.service

[Unit]
Description=Transcribe Speech to Text
After=network.target

[Service]
WorkingDirectory=/home/ahmedia/n8n/speech/  # Replace with your app's path
ExecStart=/home/ahmedia/n8n/speech/venv/bin/python3  /home/ahmedia/n8n/speech/process_speech.py  # Full path to python and app
ExecReload=/home/ahmedia/n8n/speech/venv/bin/python3  /home/ahmedia/n8n/speech/process_speech.py   
Restart=always

[Install]
WantedBy=multi-user.target
